# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:django]
command=gunicorn --bind 127.0.0.1:8000 --workers 3 Django-React-Note.wsgi:application
directory=/app
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/django.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx.log